# controller/openrc/bmaccontroller
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/usr/local/bin/bmaccontrol"
command_args="-c /etc/bmac/config.ini -d /etc/bmac/htdocs"
